# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/php-fpm.log
stderr_logfile=/var/www/storage/logs/php-fpm-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/nginx.log
stderr_logfile=/var/www/storage/logs/nginx-error.log

[program:scheduler]
command=bash -c "while true; do php /var/www/artisan schedule:run --no-interaction >> /var/www/storage/logs/scheduler.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null
